# 2do_anho/anuales/sistemas_operativos/os_laboraorio/ejercicio_facu/trabajos/entregados/configuracion/etc/xdg/autostart/disable-logout-beep.desktop
[Desktop Entry]
Type=Application
Name=Disable logout beep
Exec=xset b off
Terminal=false
OnlyShowIn=XFCE;
