# 2do_anho/anuales/sistemas_operativos/os_laboraorio/ejercicio_facu/trabajos/entregados/configuracion/etc/xdg/autostart/kali-sync-skel.desktop
[Desktop Entry]
Type=Application
Name=Copy missing config files from /etc/skel
Exec=/bin/sh -c "cp -rnT /etc/skel/ ~"
